# files/system/usr/lib/systemd/system/opt.mount
[Unit]
Description=Mount /opt

[Mount]
What=/var/opt
Where=/opt
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
